# launchd/com.katab.krc-updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.katab.krc-updater</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__KRC_HOME__/bin/krc-auto-update.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>__USER_HOME__</string>
    <key>KRC_HOME</key>
    <string>__KRC_HOME__</string>
  </dict>

  <!-- Run every hour (3600 seconds) -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>__KRC_HOME__/logs/auto-update.log</string>

  <key>StandardErrorPath</key>
  <string>__KRC_HOME__/logs/auto-update.log</string>

  <key>UserName</key>
  <string>__USER__</string>
</dict>
</plist>
